# .deployment/files/etc/systemd/system/pacman-db-sync.service
[Unit]
Description=Pacman database update

[Service]
Type=oneshot
ExecStart=/usr/bin/pacman --sync --refresh
