# .deployment/files/etc/systemd/system/pacman-db-sync.timer
[Unit]
Description=Periodic pacman sync and refresh of the database

[Timer]
OnCalendar=hourly

[Install]
WantedBy=timers.target
