# .deployment/files/etc/systemd/system/reflector.timer
[Unit]
Description=Periodic pacman mirrorlist refresh

[Timer]
OnCalendar=*-*-* 12:30:00
RandomizedDelaySec=1800

[Install]
WantedBy=timers.target
